# crontab
0 15 * * * python main.py process cbs --from_email --username $MAILUSER --password $MAILPASS
