# crontab
0,30 * * * * python main.py process_data
